# board/rockchip/common/base/etc/udev/rules.d/firefly-otg-upgrade.rules
ACTION!="add",GOTO="farsight"
KERNEL=="sd[a-z]",RUN+="/usr/local/bin/firefly-otg-control.sh start"
LABEL="farsight"
