# board/rockchip/common/base/etc/udev/rules.d/firefly-otg-upgrade.rules
ACTION!="add",GOTO="farsight"
KERNEL=="sd[a-z]",RUN+="/bin/touch /tmp/otg_start"
LABEL="farsight"
